# recipes-support/cpu-performance/cpu-performance/cpu-performance.service

[Unit]
Description=Set CPU governor to performance

[Service]
ExecStart=/usr/bin/cpuperf
Type=oneshot

[Install]
WantedBy=multi-user.target
